# api_geonature-service.conf
[program:api_geonature]
command = APP_PATH/backend/start_gunicorn.sh
autostart=true
autorestart=true
stdout_logfile = APP_PATH/var/log/api_geonature-errors.log
redirect_stderr = true