# gncitizen_frontssr-service.conf
[program:geonature]
command = APP_PATH/frontend/start_server.sh
autostart=true
autorestart=true
stdout_logfile = APP_PATH/var/log/geonature-errors.log
redirect_stderr = true